# debout/etc/udev/rules.d/00-test.rules
ACTION=="add", SUBSYSTEM=="usb", RUN+="/usr/share/debtest/test.sh"
