# deb/etc/udev/rules.d/00-test.rules
ACTION=="add", SUBSYSTEM=="usb", RUN+="/opt/debtest/test.sh"
